# sqldump2restic.timer
[Unit]
Description=scheduled SQL Database backup

[Timer]
OnCalendar=hourly
Persistent=true
#RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
